# docs/cloud/systemd_operator_api.example.service
[Unit]
Description=POVERTY_KILLER Operator API PAPER-only
After=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/poverty_killer
EnvironmentFile=/etc/poverty_killer/operator.env
ExecStart=/opt/poverty_killer/venv/bin/python -m uvicorn app.api.operator_readonly_api:create_operator_app --factory --host 127.0.0.1 --port 8765
Restart=on-failure
RestartSec=10
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
